# debian/etc/skel/.config/systemd/gnunet-userpeer.service
[Unit]
Description=GNUnet User peer

[Service]
Type=oneshot
ExecStart=/usr/bin/gnunet-arm -s

[Install]
WantedBy=default.target
